# .crontab
etup shell, email address for reporting, and environment path
SHELL=/bin/bash
MAILTO=mfreema1@stevens.edu
PATH=$HOME/biPATH=/opt/someApp/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/s$

# Example of job definition:
# .---------------- minute (0 - 59)
# |  .------------- hour (0 - 23)
# |  |  .---------- day of month (1 - 31)
# |  |  |  .------- month (1 - 12) OR jan,feb,mar,apr ...
# |  |  |  |  .---- day of week (0 - 6) (Sunday=0 or 7) OR sun,mon,tue,wed,thu,fri,sat
# |  |  |  |  |
# *  *  *  *  * user-name  command to be executed

# Every 5 min, go get a few measurements and dump them into text file
*/5 * * * * $HOME/bin/take_reading.sh

# Every 4 hours, grab that file and offload to S3, then delete it on disk
* */4 * * * $HOME/bin/upload.sh

